# ingester.timer
[Unit]
Description=Run ingester service daily

[Timer]
OnCalendar=*-*-* 4:00:00
Persistent=true

[Install]
WantedBy=timers.target
