# ingester.timer
[Unit]
Description=Run ingester service daily

[Timer]
OnCalendar=*-*-* 06:00:00

[Install]
WantedBy=timers.target
